# redis/instalacionTradicional/redis_centinela_systemd.service
[Unit]
Description=Redis Centinela

[Service]
ExecStart=/usr/local/bin/redis-server /etc/redis/centinela.conf --sentinel
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target